# cloudflare/cloud-divinity/supervisord.conf
# Supervisord config for Divinity Cloud
# Manages: Xvfb (display), Fluxbox (window manager), x11vnc (VNC server),
#          websockify (noVNC HTML5 bridge), PulseAudio (audio), Divinity (app)

[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

# 1. Xvfb — virtual framebuffer (Electron renders here)
[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 %(ENV_RESOLUTION)s -ac -nolisten tcp
priority=10
autostart=true
autorestart=true
stdout_logfile=/var/log/xvfb.log
stderr_logfile=/var/log/xvfb-error.log

# 2. Fluxbox — lightweight window manager (Electron needs one)
[program:fluxbox]
command=/usr/bin/fluxbox -display %(ENV_DISPLAY)s
priority=20
autostart=true
autorestart=true
environment=DISPLAY="%(ENV_DISPLAY)s"
stdout_logfile=/var/log/fluxbox.log
stderr_logfile=/var/log/fluxbox-error.log

# 3. x11vnc — VNC server (shares the Xvfb display)
[program:x11vnc]
command=/usr/bin/x11vnc -display %(ENV_DISPLAY)s -forever -shared -noauth -rfbport 5900
priority=30
autostart=true
autorestart=true
stdout_logfile=/var/log/x11vnc.log
stderr_logfile=/var/log/x11vnc-error.log

# 4. websockify — bridges VNC to WebSocket (noVNC HTML5 client connects here)
[program:websockify]
command=/usr/bin/websockify --web=/usr/share/novnc 6080 localhost:5900
priority=40
autostart=true
autorestart=true
stdout_logfile=/var/log/websockify.log
stderr_logfile=/var/log/websockify-error.log

# 5. PulseAudio — virtual audio (so TTS can "play" and the user hears it)
[program:pulseaudio]
command=/usr/bin/pulseaudio --system --disallow-exit --disable-shm -n -F /etc/pulse/default.pa
priority=50
autostart=true
autorestart=true
stdout_logfile=/var/log/pulseaudio.log
stderr_logfile=/var/log/pulseaudio-error.log

# 6. Divinity Electron app — the main event
# Runs the packaged Electron app with the virtual display.
# DIVINITY_WORKDIR points to the persistent volume (/data).
[program:divinity]
command=/app/apps/x/apps/main/out/Divinity-linux-arm64/divinity
directory=/app/apps/x/apps/main/out/Divinity-linux-arm64
priority=90
autostart=true
autorestart=true
startsecs=5
environment=DISPLAY="%(ENV_DISPLAY)s",DIVINITY_WORKDIR="%(ENV_DIVINITY_WORKDIR)s"
stdout_logfile=/var/log/divinity.log
stderr_logfile=/var/log/divinity-error.log
